# files/etc/service/php-fpm/run
#!/bin/bash -eu

echo "Starting php-fpm..."

if [ ! -d /run/php ]; then
	mkdir /run/php
fi

exec /usr/sbin/php-fpm7.2 --nodaemonize --fpm-config /etc/php/7.2/fpm/php-fpm.conf
