# files/etc/service/php-fpm/run
#!/bin/bash -eu

if [ ! -d /run/php ]; then
	mkdir /run/php
fi

exec /usr/sbin/php-fpm7.0 --nodaemonize --fpm-config /etc/php/7.0/fpm/php-fpm.conf
